# src/systemd/elasticdns.service
[Unit]
Description=Updates a Route53 record set with the system's current public IP

[Service]
Type=simple
ExecStart=/usr/bin/elasticdns
EnvironmentFile=/etc/elasticdns.conf

[Install]
WantedBy=multi-user.target